# deploy/systemd/carrefour-mcp.service
[Unit]
Description=carrefour-mcp MCP server
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=coursicota
Group=coursicota
WorkingDirectory=/home/coursicota/apps/carrefour-mcp/dist
EnvironmentFile=-/home/coursicota/etc/carrefour-mcp.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
UMask=0077

[Install]
WantedBy=multi-user.target
